# systemd/system/custom-startup-script.service
[Unit]
Description=Custom Startup Script

[Service]
ExecStart=/path/to/your/script.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target


